# letsencrypt-certbot/etc/crontab
41 6,15 * * * /usr/bin/certbot renew --no-self-upgrade
